# tftp/systemd/tftp.service
[Unit]
Description=Tftp Server

[Service]
ExecStart=/usr/sbin/in.tftpd -u tftp -s /var/lib/tftpboot
StandardInput=socket
