# logrotate/systemd/logrotate.service
[Unit]
Description=logrotate

[Service]
ExexStart=/usr/sbin/logrotate /etc/logrotate.d/*
