# dma/systemd/dma-cleanup-spool.service
[Unit]
Description=Cleanup of the DMA mail spool directory

[Service]
ExecStart=/usr/sbin/dma-cleanup-spool
