# dma/systemd/dma-queue-flush.timer
[Unit]
Description=Periodly flushing of the DMA mail queue.

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target
