# netplug/systemd/netplug.service
[Unit]
Description=netplug Daemon

[Service]
ExecStart=/sbin/netplugd -F

[Install]
WantedBy=multi-user.target
