# isdn4k-utils/systemd/capi.service
[Unit]
Description=Capi services for passiv ISDN cards
After=syslog.target

[Service]
ExecStart=/usr/sbin/capiinit start
ExecStop=/usr/sbin/capiinit stop

[Install]
WantedBy=multi-user.target
