# acpid/systemd/acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/acpid

[Install]
WantedBy=multi-user.target
