# openldap/systemd/openldap.service
[Unit]
Description=OpenLDAP
After=basic.target sockets.target

[Service]
Type=forking
ExecStart=/usr/sbin/slapd -u ldap -h 'ldapi:// ldaps://'
