# openldap/systemd/openldap.service
[Unit]
Description=OpenLDAP
After=basic.target

[Service]
Type=forking
ExecStart=/usr/sbin/slapd

[Install]
WantedBy=multi-user.target
