# openldap/systemd/openldap.socket
[Socket]
ListenStream=/run/ldapi

[Install]
WantedBy=sockets.target
