# openldap/systemd/openldap.service
[Unit]
Description=OpenLDAP
After=basic.target sockets.target

[Service]
ExecStart=/usr/sbin/slapd -u ldap -h 'ldapi://'
