# rng-tools/systemd/rngd.service
[Unit]
Description=Hardware RNG Entropy Gatherer Daemon

[Service]
ExecStart=/usr/sbin/rngd -f --no-tpm=1

[Install]
WantedBy=multi-user.target
